# deploy/retina-backup.service
[Unit]
Description=RETINA nightly backup to Backblaze B2
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/tower-finder
ExecStart=/usr/bin/docker compose exec -T tower-finder python3 backend/scripts/backup_to_b2.py
Environment=B2_KEY_ID=
Environment=B2_APP_KEY=
Environment=B2_BUCKET_NAME=
EnvironmentFile=-/opt/tower-finder/backend/.env
TimeoutStartSec=600
